# blinkyoureyes.service.in
[Unit]
Description=blinkyoureyes daemon service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart={{EXEFILE}}
Restart=always
RestartSec=30

[Install]
WantedBy=xsession.target
